# nyptid-studio.service
[Unit]
Description=NYPTID Studio Engine (Viral Shorts AI Platform)
After=network.target

[Service]
Type=simple
User=omatic657
WorkingDirectory=/home/omatic657/asd
Environment=FRONTEND_DIST_DIR=/home/omatic657/asd/ViralShorts-App/dist
ExecStart=/usr/bin/python3 -m uvicorn backend:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/tmp/nyptid-studio.log
StandardError=append:/tmp/nyptid-studio.log

[Install]
WantedBy=multi-user.target
